# 1_CODE_代码/okx-contract-signal-system/deployment/systemd/okx-signal-health.timer
[Unit]
Description=Periodic OKX signal health check

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true
Unit=okx-signal-health.service

[Install]
WantedBy=timers.target
